# systemd/logbook-printer.service
[Unit]
Description=Print latest logbook entries

[Service]
Type=oneshot
ExecStart=/usr/bin/node /home/pi/logbook-printer/index.mjs /dev/rfcomm0
