# systemd/logbook-printer.service
[Unit]
Description=Print latest logbook entries

[Service]
Type=oneshot
ExecStart=/usr/local/bin/node /home/pi/logbook-printer /dev/rfcomm0
